# test.nginx.service
[Unit]
Description=Тестовый юнит для ТЗ мониторинга процесса nginx

[Service]
ExecStart=/opt/testscript.sh
User=root
Group=root

